# weather-api.service
[Unit]
Description=Weather API Service
After=network.target

[Service]
WorkingDirectory=/home/deployer/weather-api/publish
ExecStart=/usr/bin/dotnet /home/deployer/weather-api/publish/weather-api.dll
Restart=always
RestartSec=5
SyslogIdentifier=weather-api
Environment=ASPNETCORE_ENVIRONMENT=Production

[Install]
WantedBy=multi-user.target
